# systemd/omarchy-check-upstream.service
[Unit]
Description=Check for upstream changes in omarchy-pkgs repository
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/omarchy-pkgs/bin/check-upstream
Environment=OMARCHY_PKGS_DIR=/root/omarchy-pkgs
Environment=OMARCHY_STATE_DIR=/root/.state

[Install]
WantedBy=multi-user.target
