# systemd/omarchy-check-versions.service
[Unit]
Description=Check PKGBUILD versions against published repo versions
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/git pull --ff-only
ExecStart=/root/omarchy-pkgs/bin/check-versions
Environment=OMARCHY_STATE_DIR=/root/.state
WorkingDirectory=/root/omarchy-pkgs

[Install]
WantedBy=multi-user.target
