# systemd/omarchy-check-versions.timer
[Unit]
Description=Check package versions every 6 hours

[Timer]
# Run every 6 hours at :30 (after GitHub sync at :00)
OnCalendar=*-*-* 00,06,12,18:30:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
