# systemd/omarchy-auto-release-edge.timer
[Unit]
Description=Auto-release edge mirror every 6 hours

[Timer]
# Run every 6 hours at :00 (after version check at :30)
OnCalendar=*-*-* 01,07,13,19:00:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
