# systemd/omarchy-auto-release-stable.service
[Unit]
Description=Auto-release omarchy-pkgs for stable mirror
After=network-online.target omarchy-check-versions.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'source /root/.omarchy/build-credentials && /root/omarchy-pkgs/bin/auto-release stable'
Environment=OMARCHY_STATE_DIR=/root/.state
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
