# pkgbuilds/wl-clip-persist/wl-clip-persist.service
[Unit]
Description = Wayland clipboard persistence daemon
ConditionEnviroment = WAYLAND_DISPLAY

[Service]
Type = simple
ExecStart = /usr/bin/wl-clip-persist --clipboard regular
Restart = on-failure

[Install]
WantedBy = default.target

