# pkgbuilds/edge/intel-ipu7-camera/camera-feed-watch.service
[Unit]
Description=Camera feed for v4l2loopback (on-demand, LED off when idle)
After=pipewire.service camera-init.service

[Service]
Type=simple
ExecStart=/usr/bin/camera-feed-watch
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
