# pkgbuilds/edge/intel-ipu7-camera/intel-ipu7-camera.service
[Unit]
Description=Start Intel IPU7 camera stack after graphical boot
After=graphical.target
ConditionPathExists=/usr/lib/systemd/system/camera-init.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl start --no-block camera-init.service v4l2-relayd@ipu7.service

[Install]
WantedBy=graphical.target
